# deploy/deliver.conf
description "Starts deliver"

start on runlevel [2345]
stop on runlevel [!2345]

chdir /folder/for/deliver
exec sudo -u sirech /path/to/python updater.py
respawn
